# keyboard-wakeup.rules
ACTION=="add", SUBSYSTEM=="input", KERNEL=="event*", RUN+="/bin/sh -c 'echo enabled > /sys/devices/virtual/input/input8/power/wakeup'"
